# project/budget_tracker.service

[Unit]
Description=Budget Tracker Web Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=budget_tracker
WorkingDirectory=/path/to/your/app  # Replace with actual path
Environment=PATH=/path/to/your/venv/bin  # Replace with actual venv path
EnvironmentFile=/path/to/your/app/.env.production  # Replace with actual path
ExecStart=/path/to/your/venv/bin/gunicorn --config gunicorn.conf.py app_optimized:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
